# python-scripts/invoice-system/com.graydavis.incometracker.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.graydavis.incometracker</string>

    <key>ProgramArguments</key>
    <array>
        <string>/usr/bin/python3</string>
        <string>/Users/graydavis28/Desktop/my-project/python-scripts/invoice-system/main.py</string>
        <string>scan-all</string>
    </array>

    <key>WorkingDirectory</key>
    <string>/Users/graydavis28/Desktop/my-project/python-scripts/invoice-system</string>

    <key>StartCalendarInterval</key>
    <dict>
        <key>Hour</key>
        <integer>8</integer>
        <key>Minute</key>
        <integer>0</integer>
    </dict>

    <key>StandardOutPath</key>
    <string>/Users/graydavis28/Desktop/my-project/python-scripts/invoice-system/tracker.log</string>

    <key>StandardErrorPath</key>
    <string>/Users/graydavis28/Desktop/my-project/python-scripts/invoice-system/tracker-error.log</string>
</dict>
</plist>
